# infrastructure/systemd/tad-agent-start.timer
[Unit]
Description=Fenêtre de démarrage des agents Codex TAD

[Timer]
OnCalendar=*-*-* 19:30:00 Africa/Abidjan
Persistent=true
Unit=tad-agent-start.service

[Install]
WantedBy=timers.target
