# infrastructure/systemd/tad-agent-report.timer
[Unit]
Description=Rapport matinal des agents Codex TAD

[Timer]
OnCalendar=*-*-* 06:05:00 Africa/Abidjan
Persistent=true
Unit=tad-agent-report.service

[Install]
WantedBy=timers.target
